# altermeshfc/altermeshfc/firmcreator/test_data/profiles/test.org.ar/include_files/etc/uci-defaults/add-opkg-repo-altermundi
#!/bin/sh
rev=$(cat /etc/openwrt_version)
line="src/gz altermundi http://openwrt.altermundi.net/snapshots/$rev/ar71xx/packages"
file="/etc/opkg.conf"

sed "\|^src/gz altermundi |d" -i $file
echo "$line" >> $file
