# altermeshfc/altermeshfc/firmcreator/test_data/profiles/test.org.ar/include_files/etc/uci-defaults/set-timezone-art3
#!/bin/sh
timezone="ART3"

uci set system.@system[0].timezone="$timezone"
uci commit system
